# devstack/supervisord.conf
[supervisord]
nodaemon=true

[program:devstack]
user=stack
command=/github/devstack/stack.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
